# com.helloworld.refresh-approved-projects.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.refresh-approved-projects</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ajiteshkumar/Documents/projects/ysws-dashboard/refresh_projects.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ajiteshkumar/Documents/projects/ysws-dashboard/refresh.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/ajiteshkumar/Documents/projects/ysws-dashboard/refresh.error</string>
</dict>
</plist>

